# docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl = unix:///var/run/supervisord.sock
username = usersupervisord
password = vuqu0shit6Ae

[program:php-fpm]
command = /usr/local/sbin/php-fpm --force-stderr --nodaemonize --pid /var/run/php-fpm.pid --fpm-config /usr/local/etc/php-fpm.conf
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
startretries=3
restartpause=10
stopsignal=QUIT

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
startretries=3
restartpause=10
stopsignal=QUIT
